# alphascreener/scheduler/systemd/alphascreener-scheduler.service
[Unit]
Description=Alpha Screener APScheduler Daemon
Documentation=https://github.com/Qanora/alpha-screener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alphascreener
Group=alphascreener
StateDirectory=alphascreener
StateDirectoryMode=0750
WorkingDirectory=/opt/alpha-screener
Environment="PYTHONUNBUFFERED=1"
Environment="ALPHASCREENER_HOME=/var/lib/alphascreener"
EnvironmentFile=/opt/alpha-screener/.env
ExecStart=/opt/alpha-screener/.venv/bin/python -m alphascreener.scheduler.run
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alphascreener-scheduler

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/alphascreener
ReadOnlyPaths=/opt/alpha-screener /opt/alpha-screener/.venv

# Resource limits
MemoryMax=4G
CPUQuota=400%
TasksMax=64
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
